# shepherd/supervisord.conf
[supervisord]
nodaemon=true

[program:logtocdx]
command=/python-shepherd/bin/logtocdx
autorestart = true
startsecs = 0