# shepherd/supervisord.conf
[supervisord]
nodaemon=true

[program:harchived]
command=/python-shepherd/agents/harchiverd.py --log-level DEBUG --amqp-url %(ENV_AMQP_URL)s --webrender-url http://webrender:8000/webtools/domimage heritrix FC-2-uris-to-render uris-to-render /images
autorestart = true
startsecs = 0
stdout_events_enabled = true
stderr_events_enabled = true

[program:uristocdxserver]
command=/python-shepherd/agents/uristocdxserver.py --amqp-url %(ENV_AMQP_URL)s --cdxserver-url %(ENV_CDXSERVER_URL)s heritrix %(ENV_QUEUE_NAME)s
autorestart = true
startsecs = 0
stdout_events_enabled = true
stderr_events_enabled = true

[program:docstow3act]
command=/python-shepherd/agents/docstow3act.py --amqp-url %(ENV_AMQP_URL)s --w3act-url %(ENV_W3ACT_URL)s --w3act-user %(ENV_W3ACT_USER)s --w3act-pw  %(ENV_W3ACT_PW)s heritrix FC-5-documents-to-catalogue
autorestart = true
startsecs = 0
stdout_events_enabled = true
stderr_events_enabled = true

[program:docstow3act]
command=/setup-testrun.sh
autorestart = false
startsecs = 0
stdout_events_enabled = true
stderr_events_enabled = true

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler
