# shepherd/supervisord.conf
[supervisord]
nodaemon=true

[program:logtocdx]
command=/python-shepherd/bin/logtocdx.py
autorestart = true
startsecs = 0